# scripts/systemd/prefect-server.service
[Unit]
Description=Prefect Server (Self-Hosted)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${VM_USER}
WorkingDirectory=/opt/fund-lens-etl
Environment="PATH=/opt/fund-lens-etl/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/fund-lens-etl/.venv/bin/prefect server start --host 0.0.0.0
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prefect-server

[Install]
WantedBy=multi-user.target
